# docker-demo.conf.upstart
description "NethServer docker demo"
author "Davide Principi <davide.principi@nethesis.it>"
respawn
respawn limit 10 5
start on started docker
stop on stopping docker

#
# To initially create the docker-demo container:
#
# docker create -h demo.nethserver.org --name docker-demo \
#   -p 8080:980 nethserver/docker-demo:latest
#

exec docker start -a docker-demo
